# addons/vivobook/tdp-config.service
[Unit]
Description=TDP Vivobook Config

[Service]
ExecStart=/usr/bin/low-tdp-ryzen

[Install]
WantedBy=multi-user.target
